# package/elsewhere/DEBIAN/postinst
#!/bin/sh

systemctl enable elsewhere-splash.service
