# services/splash.service
[Unit]
Description=boot splash screen

[Service]
ExecStart=/usr/bin/fim /home/pi/splash.png

[Install]
WantedBy=basic.target
