# package/elsewhere/etc/systemd/system/elsewhere-splash.service
[Unit]
Description=Elsewhere splash screen
After=local-fs.target

[Service]
ExecStart=/bin/bash /usr/share/elsewhere/splash/set.sh

[Install]
WantedBy=sysinit.target
